# kloxo/httpdocs/theme/filecore/init.crontab
*   *   *   *   *   sh /usr/local/lxlabs/kloxo/bin/fetchmailticket.sh
