# kloxo/init/kloxo-wrapper.cron
* * * * * root sh /script/run-kloxo-wrapper &